# daemons/retry-queue-ane2.timer
[Unit]
Description=Timer Retry Queue Service ANE2

[Timer]
OnBootSec=1min
OnUnitInactiveSec=300s
AccuracySec=1s

[Install]
WantedBy=timers.target
