# daemons/status-ane2.timer
[Unit]
Description=Timer Status Service ANE2

[Timer]
OnBootSec=1min
OnUnitInactiveSec=30s
AccuracySec=1s

[Install]
WantedBy=timers.target
